# scripts/systemd/hashpool-stats-proxy.service
[Unit]
Description=Hashpool Proxy Statistics Service
After=network.target

[Service]
Type=simple
User=hashpool
Group=hashpool
WorkingDirectory=/opt/hashpool
ExecStartPre=+/bin/sh -c 'mkdir -p /var/log/hashpool /var/lib/hashpool/stats-proxy && chown hashpool:hashpool /var/log/hashpool /var/lib/hashpool && chmod 755 /var/log/hashpool /var/lib/hashpool /var/lib/hashpool/stats-proxy'
ExecStart=/opt/hashpool/bin/stats_proxy \
  --config /opt/hashpool/config/stats-proxy.config.toml \
  --tproxy-config /opt/hashpool/config/tproxy.config.toml \
  --shared-config /opt/hashpool/config/shared/miner.toml \
  --db-path /var/lib/hashpool/stats-proxy/stats.db \
  -f /var/log/hashpool/stats-proxy.log
StandardOutput=file:/var/log/hashpool/stats-proxy.log
StandardError=file:/var/log/hashpool/stats-proxy.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
