# scripts/systemd/hashpool-web-pool.service
[Unit]
Description=Hashpool Pool Web Dashboard
After=network.target hashpool-stats-pool.service
Requires=hashpool-stats-pool.service

[Service]
Type=simple
User=hashpool
Group=hashpool
WorkingDirectory=/opt/hashpool
ExecStartPre=+/bin/sh -c 'mkdir -p /var/log/hashpool && touch /var/log/hashpool/web-pool.log && chown hashpool:hashpool /var/log/hashpool /var/log/hashpool/web-pool.log && chmod 755 /var/log/hashpool && chmod 644 /var/log/hashpool/web-pool.log'
ExecStart=/opt/hashpool/bin/web_pool \
  --web-pool-config /opt/hashpool/config/web-pool.config.toml \
  --shared-config /opt/hashpool/config/shared/pool.toml \
  -f /var/log/hashpool/web-pool.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
