# scripts/systemd/hashpool-bitcoind.service
[Unit]
Description=Bitcoin daemon (testnet4)
After=network.target

[Service]
Type=simple
User=hashpool
Group=hashpool
WorkingDirectory=/opt/hashpool
ExecStartPre=+/bin/sh -c 'mkdir -p /var/log/hashpool /var/lib/hashpool/bitcoind/testnet4 && touch /var/log/hashpool/bitcoind.log && chown -R hashpool:hashpool /var/log/hashpool /var/lib/hashpool && chmod 755 /var/log/hashpool /var/lib/hashpool /var/lib/hashpool/bitcoind /var/lib/hashpool/bitcoind/testnet4 && chmod 644 /var/log/hashpool/bitcoind.log'
ExecStart=/opt/hashpool/bin/bitcoind \
  -datadir=/var/lib/hashpool/bitcoind \
  -chain=testnet4 \
  -conf=/opt/hashpool/config/bitcoin.conf \
  -nodaemon
ExecStop=/opt/hashpool/bin/bitcoin-cli -datadir=/var/lib/hashpool/bitcoind -chain=testnet4 stop
Restart=on-failure
RestartSec=10
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
