# scripts/systemd/hashpool-stats-pool.service
[Unit]
Description=Hashpool Pool Statistics Service
After=network.target

[Service]
Type=simple
User=hashpool
Group=hashpool
WorkingDirectory=/opt/hashpool
ExecStartPre=+/bin/sh -c 'mkdir -p /var/log/hashpool /var/lib/hashpool/stats-pool && touch /var/log/hashpool/stats-pool.log && chown hashpool:hashpool /var/log/hashpool /var/log/hashpool/stats-pool.log /var/lib/hashpool && chmod 755 /var/log/hashpool /var/lib/hashpool /var/lib/hashpool/stats-pool && chmod 644 /var/log/hashpool/stats-pool.log'
ExecStart=/opt/hashpool/bin/stats_pool \
  --tcp-address 127.0.0.1:9083 \
  --http-address 127.0.0.1:9084 \
  --config /opt/hashpool/config/prod/stats-pool.config.toml \
  --metrics-db-path /var/lib/hashpool/stats-pool/metrics.db \
  -f /var/log/hashpool/stats-pool.log
StandardOutput=file:/var/log/hashpool/stats-pool.log
StandardError=file:/var/log/hashpool/stats-pool.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
