# scripts/systemd/hashpool-jd-server.service
[Unit]
Description=Hashpool Job Declarator Server
After=network.target hashpool-pool.service
Wants=hashpool-pool.service

[Service]
Type=simple
User=hashpool
Group=hashpool
WorkingDirectory=/opt/hashpool
ExecStartPre=+/bin/sh -c 'mkdir -p /var/log/hashpool && touch /var/log/hashpool/jd-server.log && chown hashpool:hashpool /var/log/hashpool /var/log/hashpool/jd-server.log && chmod 755 /var/log/hashpool && chmod 644 /var/log/hashpool/jd-server.log'
ExecStart=/opt/hashpool/bin/jd_server \
  -c /opt/hashpool/config/jds.config.toml \
  -f /var/log/hashpool/jd-server.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
