# scripts/systemd/hashpool-pool.service
[Unit]
Description=Hashpool Pool (Stratum V2)
After=network.target hashpool-stats-pool.service hashpool-mint.service
Wants=hashpool-stats-pool.service hashpool-mint.service

[Service]
Type=simple
User=hashpool
Group=hashpool
WorkingDirectory=/opt/hashpool
ExecStart=/opt/hashpool/bin/pool_sv2 \
  -c /opt/hashpool/config/pool.config.toml \
  -g /opt/hashpool/config/shared/pool.toml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hashpool-pool

[Install]
WantedBy=multi-user.target
