# scripts/systemd/hashpool-proxy.service
[Unit]
Description=Hashpool Proxy (Translator SV2)
After=network.target hashpool-stats-proxy.service hashpool-pool.service
Requires=hashpool-stats-proxy.service

[Service]
Type=simple
User=hashpool
Group=hashpool
WorkingDirectory=/opt/hashpool
Environment="BITCOIND_NETWORK=testnet4"
Environment="RUST_LOG=info"
Environment="CDK_WALLET_DB_PATH=/var/lib/hashpool/translator/wallet.sqlite"
ExecStartPre=+/bin/sh -c 'mkdir -p /var/log/hashpool /var/lib/hashpool/translator && touch /var/log/hashpool/proxy.log && chown hashpool:hashpool /var/log/hashpool /var/log/hashpool/proxy.log /var/lib/hashpool && chmod 755 /var/log/hashpool /var/lib/hashpool /var/lib/hashpool/translator && chmod 644 /var/log/hashpool/proxy.log'
ExecStart=/opt/hashpool/bin/translator_sv2 \
  -c /opt/hashpool/config/tproxy.config.toml \
  -g /opt/hashpool/config/shared/miner.toml \
  -f /var/log/hashpool/proxy.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
